# scripts/knowledge-base/com.user.transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Auto-transcription LaunchAgent for macOS

  This plist watches ~/audio-inbox/ for new files and triggers transcription.

  Installation:
    1. Copy to ~/Library/LaunchAgents/com.user.transcribe.plist
    2. Replace YOUR_USERNAME with your actual username
    3. Replace hf_your_token_here with your Hugging Face token
    4. Load with: launchctl load ~/Library/LaunchAgents/com.user.transcribe.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.user.transcribe.plist
    rm ~/Library/LaunchAgents/com.user.transcribe.plist

  Logs:
    tail -f /tmp/transcribe.log
    tail -f /tmp/transcribe.err
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.transcribe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/scripts/transcribe.sh</string>
    </array>

    <!-- Watch for new files in audio-inbox -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USERNAME/audio-inbox</string>
    </array>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HF_TOKEN</key>
        <string>hf_your_token_here</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/transcribe.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/transcribe.err</string>

    <!-- Keep alive (restart if crashes) -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run at load (optional) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
