# scripts/knowledge-base/com.user.transcript-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Transcript Watcher LaunchAgent for macOS

  This plist runs the transcript-watcher.sh script as a persistent background service.
  The watcher uses fswatch to monitor for new transcripts and processes them via openclaw agent.

  Installation:
    1. Copy to ~/Library/LaunchAgents/com.user.transcript-watcher.plist
    2. Replace YOUR_CHAT_ID_HERE with your Telegram chat ID (from @userinfobot)
    3. Load with: launchctl load ~/Library/LaunchAgents/com.user.transcript-watcher.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.user.transcript-watcher.plist
    rm ~/Library/LaunchAgents/com.user.transcript-watcher.plist

  Logs:
    tail -f /tmp/transcript-watcher.log
    tail -f /tmp/transcript-watcher.err

  Verify:
    launchctl list | grep transcript-watcher
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.transcript-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/My Shared Files/scripts/transcript-watcher.sh</string>
    </array>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>YOUR_CHAT_ID_HERE</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Volumes/My Shared Files</string>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/transcript-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/transcript-watcher.err</string>

    <!-- Keep alive - restart if crashes (this is a long-running fswatch loop) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Run at load - start watching immediately on boot/login -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
